# scripts/gemma-cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leonard.gemma.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/leonard/.claude/scripts/gemma-cron-daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/leonard/.claude/cache/gemma-cron.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/leonard/.claude/cache/gemma-cron.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
